# auto_start_program/auto_start_app.service
[Unit]
Description=BlaBlaBla

[Service]
Type=idle
ExecStart=/home/expert/auto_start/auto_start.sh

[Install]
WantedBy=multi-user.target
